# teleport.conf
# teleport

description "An simple service to manage lots of servers."
# author "Inumedia <inumedia@inumedia.net>"

start on runlevel [2345]
respawn
script
exec teleport start
end script